# docker-compose/crontab
* * * * * su www-data -s /bin/sh -c "/var/www/artisan schedule:run" >> /dev/stdout 2>&1

